# docker-app-base/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=service nginx start

[program:puma]
command=/usr/local/bin/pumarun

[program:mysqld]
command=/usr/local/bin/run
